# sss.conf
[program:sss]
command=python3 /home/curator2/sss/main.py
autostart=true
autorestart=true
stderr_logfile=/home/curator2/sss/test.err.log
stdout_logfile=/home/curator2/sss/test.out.log